# playbooks/roles/sslh/templates/sslh.service.j2
[Service]
PrivateTmp=true
StandardOutput=null
RestartSec=5s
Restart=on-failure
PIDFile={{ sslh_pid_file }}
