# playbooks/roles/wireguard/templates/streisand-wireguard.sh
#!/bin/sh
### BEGIN INIT INFO
# Provides:          streisand-wireguard
# Required-Start:    $network $remote_fs $local_fs
# Required-Stop:     $network $remote_fs $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Persist WireGuard firewall rules for Streisand
### END INIT INFO

{% for rule in wireguard_firewall_rules %}
{{ rule }}
{% endfor %}

exit 0
