# playbooks/roles/nginx/templates/nginx.service.j2
[Service]
PrivateTmp=true
RestartSec=5s
Restart=on-failure
