# playbooks/roles/ad-blocking/files/download-blocklists.service
[Unit]
Description=Download blocklists

[Service]
Type=oneshot
ExecStart=/usr/local/bin/download-blocklists

