# playbooks/roles/shadowsocks/templates/shadowsocks-libev.service.j2
[Service]
PrivateTmp=true
Restart=on-failure
RestartSec=5s
