# playbooks/roles/openconnect/templates/streisand-openconnect-service.sh.j2
#!/bin/sh
### BEGIN INIT INFO
# Provides:          streisand-openconnect
# Required-Start:    $network $remote_fs $local_fs
# Required-Stop:     $network $remote_fs $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Persist OpenConnect firewall rules for Streisand
### END INIT INFO

{{ ocserv_firewall_rule }}
